# src/systemd/tuptime.timer
[Unit]
Description=Tuptime service timer

[Timer]
OnBootSec=1min
OnCalendar=*:0/3

[Install]
WantedBy=timers.target
