# src/systemd/tuptime.timer
[Unit]
Description=Tuptime service timer

[Timer]
OnBootSec=1min
OnActiveSec=3min

[Install]
WantedBy=timers.target
