# src/rc.d/openbsd/tuptime
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/local/bin/tuptime"
daemon_flags="-x"

rc_reload=NO
rc_stop=NO
rc_check=NO

. /etc/rc.d/rc.subr

rc_cmd $1
