# deb-files/package-tuptime-deb-2.2.00/DEBIAN/postinst
#!/bin/sh
set -e


if [ -f "/etc/systemd/system/tuptime.service" ]; then
	systemctl enable tuptime.service
	systemctl status tuptime.service
fi

