# latest/cron.d/tuptime
# /etc/cron.d/tuptime: crontab entrie for tuptime update. 

# Execute tuptime each hour for avoid loss time if system don't shutdown correctly.
# NOTE: Decrease the execution time for increase the accuracity.

MAILTO=root

* *    * * *   root    if [ -x /usr/bin/tuptime ]; then /usr/bin/tuptime -x > /dev/null; fi

