# deb-files/package-tuptime-deb-2.6.10/DEBIAN/postinst
#!/bin/sh
set -e


if [ -f "/lib/systemd/system/tuptime.service" ]; then
        systemctl enable tuptime.service
        systemctl start tuptime.service
fi
