# src/cron.d/tuptime
# /etc/cron.d/tuptime: crontab entrie for tuptime update. 

# NOTE: Decrease the execution time for increase accuracity.


*/5 * * * *   _tuptime    if [ -x /usr/bin/tuptime ]; then /usr/bin/tuptime -x > /dev/null; fi

