# src/cron.d/tuptime
# /etc/cron.d/tuptime: crontab entry for tuptime update. 

# NOTE: Decrease the execution time for increase accuracity.


*/5 * * * *   _tuptime    [ -x /usr/bin/tuptime ] && /usr/bin/tuptime -q

