# latest/init.d/tuptime.init.d-debian7
#!/bin/sh
# tuptime - Report the historical and statistical running time of the system, keeping it between restarts.
# Copyright (C) 2015 - Ricardo F.

# This program is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, either version 2 of the License, or
# (at your option) any later version.

# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
# You should have received a copy of the GNU General Public License
# along with this program.  If not, see <http://www.gnu.org/licenses/>.
#

### BEGIN INIT INFO
# Provides:          tuptime
# Required-Start:    $local_fs $time
# Required-Stop:     $local_fs $time
# Should-Start:
# Default-Start:     1 2 3 4 5
# Default-Stop:      0 6
# Short-Description: Start tuptime at boot time
# Description:       Update tuptime.
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin
. /lib/init/vars.sh
. /lib/lsb/init-functions

F_BIN='tuptime'
SERVC="/usr/bin/$F_BIN"


do_start () {
        # Start service
        log_action_begin_msg "Starting $F_BIN"
        log_end_msg 0
        $SERVC -x > /dev/null
}

do_stop () {
        # Stop service
        log_action_begin_msg "Stopping $F_BIN"
        log_end_msg 0
        $SERVC -xg > /dev/null
}

do_status () {
	# Status service
	$SERVC
}


case "$1" in
  start|"")
        do_start
        ;;
  restart)
        log_action_begin_msg "Restarting $F_BIN"
        log_end_msg 0
        do_stop
        do_start
        ;;
  stop)
        do_stop
        ;;
  status)
        do_status
        exit $?
        ;;
  *)
        echo "Usage: $F_BIN [start|stop|restart|status]" >&2
        exit 3
        ;;
esac

