# src/openrc/tuptime
#!/sbin/openrc-run

name='tuptime'
command='/usr/bin/tuptime'
command_user='tuptime'
command_args='-x'

depend() {
	need localmount
	after bootmisc
	after clock
}

stop() {
	ebegin "Stopping $RC_SVCNAME"
	command_args='-xg'
	start-stop-daemon -u "${command_user}" -exec "${command}" -- "${command_args}"
	eend $?
}
