# src/openrc/tuptime
#!/sbin/openrc-run

name='Tuptime'
command='/usr/bin/tuptime'
command_user='tuptime'

depend() {
	need localmount
	after bootmisc
	after clock
}

start() {
	ebegin "Starting $RC_SVCNAME"
	command_args='-x'
	eend $?
}
  
stop() {
	ebegin "Stopping $RC_SVCNAME"
	command_args='-xg'
	eend $?
}
