# src/openrc/tuptime
#!/sbin/openrc-run

name='tuptime'
command='/usr/bin/tuptime'
command_user='_tuptime'
command_args='-x'

depend() {
	need localmount
	after bootmisc
	after clock
}

start() {
        ebegin "Starting ${name}"
        su ${command_user} -s /bin/sh -c "${command} -x"
        eend $?
}

stop() {
        ebegin "Stopping ${name}"
        su ${command_user} -s /bin/sh -c "${command} -xg"
        eend $?
}

status() {
        ebegin "Status ${name}"
        su ${command_user} -s /bin/sh -c "${command} -x"
        eend $?
}
